# dnsmasq.service.conf
[program:dnsmasq]
command=/usr/sbin/dnsmasq -d
