# squid.service.conf
[program:squid]
command=/usr/local/bin/start.sh
stdout_events_enabled = true
stderr_events_enabled = true
