# services/dnsmasq.service.conf
[program:dnsmasq]
priority=1
command=/usr/sbin/dnsmasq -d
