# services/dnsmasq.service.conf
[program:dnsmasq]
priority=1
stopwaitsecs=5
command=/usr/sbin/dnsmasq -d
stdout_events_enabled = true
stderr_events_enabled = true
