# services/iptables.service.conf
[program:iptables]
command=/usr/local/bin/iptables.py
priority=3
stopwaitsecs=5
stdout_events_enabled = true
stderr_events_enabled = true
