# cookbooks/squid/files/default/squid.service.conf
[program:squid]
user=squid
command=squid -N
priority=2
stopwaitsecs=2
stdout_events_enabled = true
stderr_events_enabled = true
