# cookbooks/dnsmasq/files/default/dnsmasq.service.conf
[program:dnsmasq]
priority=1
stopwaitsecs=2
command=/usr/sbin/dnsmasq -d
stdout_events_enabled = true
stderr_events_enabled = true
