# cookbooks/iptables/files/default/iptables.service.conf
[program:iptables]
user=root
command=/usr/local/bin/iptables.py
priority=3
stopwaitsecs=2
stdout_events_enabled = true
stderr_events_enabled = true
